# BacklightService.service
[Unit]
Description=BacklightService

[Service]
ExecStart=/opt/BacklightService/BacklightService.py

[Install]
WantedBy=multi-user.target
